# .config/systemd/user/xsession.target
[Unit]
Description=xsession has been started
